# deploy/systemd/solana-wallet-investigator-watch.service
[Unit]
Description=Solana wallet utility-project watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/solana-wallet-investigator
EnvironmentFile=-/root/solana-wallet-investigator/.env
ExecStart=/root/solana-wallet-investigator/venv/bin/python /root/solana-wallet-investigator/auto_utility_scan.py --watch
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
